# config/.config/aconfmgr/files/etc/systemd/system/grub-btrfsd.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/grub-btrfsd --syslog --timeshift-auto
