# config/.config/aconfmgr/files/etc/systemd/system/ollama-webui.service
# container-3917cb7fe9d997344fd0083bd4bd3d62cbac493bbd9a428ded5f9070f0216e07.service
# autogenerated by Podman 5.5.2
# Thu Jul 17 19:38:49 CST 2025

[Unit]
Description=Webui of ollama by docker
Documentation=man:podman-generate-systemd(1)
Wants=network-online.target
After=network-online.target
RequiresMountsFor=/run/containers/storage

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70
ExecStart=/usr/bin/podman start 3917cb7fe9d997344fd0083bd4bd3d62cbac493bbd9a428ded5f9070f0216e07
ExecStop=/usr/bin/podman stop  \
        -t 10 3917cb7fe9d997344fd0083bd4bd3d62cbac493bbd9a428ded5f9070f0216e07
ExecStopPost=/usr/bin/podman stop  \
        -t 10 3917cb7fe9d997344fd0083bd4bd3d62cbac493bbd9a428ded5f9070f0216e07
PIDFile=/run/containers/storage/overlay-containers/3917cb7fe9d997344fd0083bd4bd3d62cbac493bbd9a428ded5f9070f0216e07/userdata/conmon.pid
Type=forking

[Install]
WantedBy=default.target
